# simlog.service
[Unit]
Description=SOC attack log simulator -> Splunk HEC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 実行ユーザ。Splunkサーバに既にある splunk アカウントを使うのが手軽
# （プロセスはHECへHTTP送信するだけなので任意のユーザでよい）
User=splunk
Group=splunk
# リポジトリ一式を置いた場所
WorkingDirectory=/opt/soc-log-simulator
ExecStart=/usr/bin/python3 /opt/soc-log-simulator/run.py
Restart=on-failure
RestartSec=5

# HECトークンは WorkingDirectory 内の config.local.json から自動的に読まれる。
# 環境変数で渡したい場合は、トークンを書いた専用ファイルを指定する（推奨 chmod 600）:
# EnvironmentFile=/opt/soc-log-simulator/simulator.env   # 中身: SPLUNK_HEC_TOKEN=xxxxxxxx

[Install]
WantedBy=multi-user.target
